# deploy/launchd/com.foreignmagazine.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.foreignmagazine.daily</string>

    <key>WorkingDirectory</key>
    <string>/Users/BGS/Documents/foreign magazine grab daily</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>set -euo pipefail; cd "/Users/BGS/Documents/foreign magazine grab daily"; if [[ -d .venv ]]; then source .venv/bin/activate; fi; exec python3 app.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PORT</key>
      <string>5050</string>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/BGS/Documents/foreign magazine grab daily/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/BGS/Documents/foreign magazine grab daily/logs/launchd.err.log</string>
  </dict>
</plist>
